# data/supervisor/default.conf
[supervisord]
nodaemon = true
logfile = /var/log/supervisor/supervisord.log
childlogdir = /var/log/supervisor
pidfile = /tmp/supervisord.pid
strip_ansi = true

[unix_http_server]
file = /var/run/supervisor.sock
chmod = 0700

[supervisorctl]
serverurl = unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:backend]
directory=/data
command=/data/services/backend/backend
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/data/logs/backend.log
; environment=ENV=production

[program:mensa]
directory=/data
command=/data/services/mensa/mensa
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/data/logs/mensa.log
environment=ENV=production

[program:svc]
directory=/data
command=/data/services/svc/svc
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/data/logs/svc.log
environment=ENV=production
